# bin/crontab/backup-system
#!/bin/env sh

set -x
set -e

backup_source="/"
backup_destination="/run/media/bryant/tiagovault/backups"

# Define exclude folders and additional excludes
exclude_folders="
  /home/bryant/desktop /home/bryant/documents /home/bryant/downloads
  /home/bryant/MEGAsync /home/bryant/pictures /home/bryant/Calibre Library
  /home/bryant/music /home/bryant/zsh /home/bryant/videos
  /home/bryant/personaldisk /home/bryant/.cache /home/bryant/.tealdeer-cache
  /home/bryant/.local/share/
"

additional_excludes="
  --exclude=/dev --exclude=/proc --exclude=/sys --exclude=/tmp
  --exclude=/run --exclude=/mnt --exclude=/media --exclude=/lost+found
  --exclude=/var/log --exclude=/var/tmp --exclude=/var/run
  --exclude=/swapfile --exclude=*/Cache --exclude=*/cache --exclude=*/.cache
"

create_directories() {
  previous_backup="$1"
  backup_folder="$2"

  mkdir -p "$previous_backup" || {
    echo "Failed to create directory: $previous_backup"
    exit 1
  }

  mkdir -p "$backup_folder" || {
    echo "Failed to create directory: $backup_folder"
    exit 1
  }
}

construct_excludes() {
  excludes=""

  for folder in $exclude_folders; do
    excludes="$excludes --exclude=$folder"
  done

  excludes="$excludes $additional_excludes"
  echo "$excludes"
}

perform_backup() {
  previous_backup="$1"
  backup_folder="$2"
  excludes="$3"

  error_log=~/.rsync_last_error.log

  rsync -aAXv --delete --progress $excludes \
    --link-dest="$previous_backup" "$backup_source" "$backup_folder" 2>"$error_log"

  if [ $? -ne 0 ]; then
    echo "Rsync encountered an error. Check $error_log for details."
  fi
}

sync_data() {
  echo "Syncing data..."
  sync
  echo "Data synced successfully."
}

# Main function to coordinate the backup process
main() {
  backup_date=$(date +%Y-%m-%d)
  previous_backup="$backup_destination/arch/previous"
  backup_folder="$backup_destination/arch/$backup_date"

  create_directories "$previous_backup" "$backup_folder"

  excludes=$(construct_excludes)

  perform_backup "$previous_backup" "$backup_folder" "$excludes"

  sync_data

  echo "Backup finalized successfully."

  exit 0
}

main
