# bin/crontab/backup-books
#!/bin/env sh

set -x
set -e

backup_source="$HOME/documents/books"
backup_destination="/run/media/bryant/tiagovault/backups/books"

backup() {
  backup_date=$(date +%Y-%m-%d)
  backup_folder="$backup_destination/$backup_date"

  mkdir -p "$backup_folder"

  rsync -aAXv --delete --progress "$1" "$backup_source/" "$backup_folder"

  if [ $? -ne 0 ]; then
    dunstify "Rsync failed" "Backup not completed"
    exit 1
  fi

  sync

  if [ $? -ne 0 ]; then
    dunstify "Sync failed" "Some data might not be written to the disk"
    exit 1
  fi

  dunstify "Backup books" "finalized successfully"
  send-sound-alert
}

get_exclude_folders() {
  exclude_folders=""

  exclude_params=""
  for folder in $exclude_folders; do
    exclude_params="$exclude_params --exclude=$folder"
  done

  echo "$exclude_params"
}

main() {
  exclude_params=$(get_exclude_folders)
  backup "$exclude_params"
}

main
