# bin/crontab/backup-notes
#!/usr/bin/env sh

NOTE_REPO="$HOME/Documents/github/notes"
YOUTUBE_REPO="$HOME/Documents/github/youtube"
DELTOCHONOTES="$HOME/Documents/github/deltocho-notes"

is_git_repo() {
	git -C "$1" rev-parse --is-inside-work-tree >/dev/null 2>&1
}

get_current_branch() {
	git -C "$1" rev-parse --abbrev-ref HEAD
}

remote_branch_exists() {
	git -C "$1" rev-parse --verify "origin/$2" >/dev/null 2>&1
}

has_unpushed_commits() {
	[ -n "$(git -C "$1" rev-list "origin/$2"..HEAD 2>/dev/null)" ]
}

has_local_changes() {
	git -C "$1" status --porcelain | grep -q .
}

push_to_remote() {
	git -C "$1" push origin "$2" >/dev/null 2>&1
}

stage_all_changes() {
	git -C "$1" add -A
}

commit_changes() {
	git -C "$1" commit -m "vault backup: $(date '+%Y-%m-%d %H:%M:%S')" >/dev/null 2>&1
}

push_existing_commits() {
	repo="$1"
	branch=$(get_current_branch "$repo")
	remote_branch_exists "$repo" "$branch" || return 0
	has_unpushed_commits "$repo" "$branch" || return 0
	push_to_remote "$repo" "$branch"
}

commit_and_push() {
	repo="$1"
	has_local_changes "$repo" || return 0
	stage_all_changes "$repo"
	commit_changes "$repo"
	push_to_remote "$repo" "$(get_current_branch "$repo")"
}

sync_repo() {
	repo="$1"
	[ -d "$repo" ] || return 1
	is_git_repo "$repo" || return 1
	push_existing_commits "$repo" && commit_and_push "$repo"
}

main() {
	sync_repo "$NOTE_REPO"
	sync_repo "$YOUTUBE_REPO"
	sync_repo "$DELTOCHONOTES"
}

main
